# templates/apache.supervisor.conf.erb
[program:apache2]
command=/bin/bash -c "mkdir -p /var/run/apache2 && source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
auto_start=true
autorestart=true