# templates/tideways.supervisor.conf.erb
[program:tideways-daemon]
command=/usr/bin/tideways-daemon --log=/var/log/tideways/daemon.log --pidfile=/var/run/tideways/tideways-daemon.pid
user=tideways
auto_start=true
autorestart=true